# data/crontab
0 23 * * * /bin/bash -c DISPLAY=:0 && /usr/bin/xset s blank
0 6 * * * /bin/bash -c DISPLAY=:0 && /usr/bin/xset s noblank

